# echotune.service
[Unit]
Description=EchoTune AI - Music Recommendation Service
Documentation=https://github.com/primoscope/doved
After=network.target mongodb.service postgresql.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=echotune
Group=echotune
WorkingDirectory=/opt/echotune
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=/opt/echotune/.env
ExecStart=/usr/bin/node src/index.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=echotune

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/echotune/logs /opt/echotune/uploads /tmp
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=1G
CPUQuota=100%

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target